# mint-dictate.service
[Unit]
Description=Mint Dictate Tray
After=graphical-session.target
Wants=graphical-session.target

[Service]
Type=simple
WorkingDirectory=%h/projects/mint-dictate
Environment=DISPLAY=:0
Environment=XAUTHORITY=%h/.Xauthority
Environment=PYTHONUNBUFFERED=1
ExecStartPre=/bin/sleep 8
ExecStart=/usr/bin/python3 %h/projects/mint-dictate/mint_dictate.py
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target
