# packaging/mint-dictate.service
[Unit]
Description=Mint Dictate Tray
After=graphical-session.target
Wants=graphical-session.target

[Service]
Type=simple
WorkingDirectory=/opt/mint-dictate
Environment=PYTHONUNBUFFERED=1
ExecStartPre=/bin/sleep 8
ExecStart=/usr/bin/mint-dictate
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target
